# infra/supervisor/webhook-relay-horizon.conf
; Supervisord config for the production Horizon worker.
; Deploys to /etc/supervisor/conf.d/webhook-relay-horizon.conf on the EC2 box.
; Reload with:  sudo supervisorctl reread && sudo supervisorctl update

[program:webhook-relay-horizon]
process_name=%(program_name)s
command=php /var/www/webhook-relay/current/artisan horizon
autostart=true
autorestart=true
user=ubuntu
redirect_stderr=true
stdout_logfile=/var/log/webhook-relay/horizon.log
stopwaitsecs=3600
